# deploy/kabanek-crm.service
[Unit]
Description=Kabanek CRM - Klienci
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/crm-gospodarstwo/kabanek-crm
ExecStart=/opt/crm-gospodarstwo/venv/bin/python app.py
Restart=always
RestartSec=5
Environment=SECRET_KEY=ZMIEN-NA-LOSOWY-KLUCZ-KABANEK

[Install]
WantedBy=multi-user.target
