# deploy/receptury-crm.service
[Unit]
Description=Kabanek CRM - Receptury
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/crm-gospodarstwo/receptury-crm
ExecStart=/opt/crm-gospodarstwo/venv/bin/python app.py
Restart=always
RestartSec=5
Environment=SECRET_KEY=ZMIEN-NA-LOSOWY-KLUCZ-RECEPTURY

[Install]
WantedBy=multi-user.target
